# config/systemd/g9-trigger-report.service
[Unit]
Description=Gate9 Trigger-based Daily Report
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/root
ExecStart=/usr/bin/docker exec gate9_sports-web bin/rails report:cycle
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
